# rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/usr/bin/with-contenv sh

# run (from /usr/local/bin)
exec qbit_manage
