# rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/usr/bin/with-contenv sh

# set CRON using docker env
echo "${CRON:-"*/15 * * * *"} /usr/local/bin/qbit_manage --run" | tee -a /etc/crontabs/root

# run crontab on foreground
exec crond -f
