# rootfs/etc/s6-overlay/s6-rc.d/qbit_manage/finish
#!/usr/bin/execlineb -P

# this will stop the whole container
/run/s6/basedir/bin/halt
